# xserver.service
[Unit]
Description=X Server

[Service]
ExecStart=startx

[Install]
WantedBy=multi-user.target